# Router_Old/etc/crontabs/root
# Update Software Lists
@weekly opkg update
#
